# print-server/prague-printer.service
[Unit]
Description=Prague Thermal Print Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/prague/print-server
ExecStart=/home/pi/prague/print-server/.venv/bin/python print_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
